# chatgpt.service
[Unit]
Description=chatgpt daemon

[Service]
ExecStart=/path-to-your-program
Restart=on-failure

[Install]
WantedBy=multi-user.target